# state_space/formal_state_space_spec_v0.json
{
  "kind": "state_space.formal_spec.v0",
  "title": "State Space Formal Spec",
  "objective": {
    "owner": "spawn",
    "control_plane": "state_space",
    "outcomes": [
      "model declared userland and observed runtime state",
      "compute dependency graphs and transition gate decisions",
      "provide deterministic artifacts for PLAN, IMPLEMENT, and LOOP"
    ],
    "provenance": {
      "derived_from": "prior session capture normalized and pruned",
      "companion_research_packet": [
        "../state_space.docs/README.md",
        "../state_space.docs/review.normalized.json",
        "../state_space.docs/high_signal_candidates.json"
      ]
    }
  },
  "scope": {
    "in_scope": [
      "state-space meta-skill design (build, update, gate, verify)",
      "runtime_surface, userland_spec, and state_space_snapshot contract shape",
      "dependency graph artifact shape and lifecycle",
      "recompute cadence (periodic + post-change)",
      "pre-commit/CI/post-commit gate integration"
    ],
    "out_of_scope": [
      "full implementation of collectors/parsers",
      "refactor of existing spawn internals beyond contract boundaries",
      "policy catalog details per service/tool"
    ]
  },
  "architecture_authority": {
    "single_authority": "spawn",
    "spawn_owns": [
      "known_state",
      "dependency_graph",
      "transition_gate_decisions"
    ],
    "sweep_boundary": "spawnctl state snapshot",
    "plan_and_implement": "consume immutable snapshot references and do not own gate-policy logic",
    "loop": "orchestration-only; no duplicated policy evaluation"
  },
  "canonical_data_model": {
    "declared_observed_state_model": {
      "runtime_surface": "declared authoritative unit for one managed app/service/binary",
      "userland_spec": "derived aggregate of runtime_surface[] for current host/profile",
      "state_space_snapshot": "observed live runtime collection",
      "rule": "missing file-declared settings are represented as empty config artifact sets, not absent contracts"
    },
    "minimum_shape": {
      "meta": [
        "schema_version",
        "snapshot_id",
        "host",
        "profile",
        "collected_at",
        "ttl_s"
      ],
      "userland_spec": [
        "runtime_surfaces_declared[]",
        "xdg_policy",
        "profiles",
        "observability_declared"
      ],
      "runtime_surface": [
        "app_id",
        "binary",
        "entrypoints",
        "launch_surfaces",
        "xdg_dirs",
        "owned_paths",
        "writable_paths",
        "env_contract",
        "systemd_units",
        "credentials",
        "ipc",
        "ports_or_sockets",
        "verification_rules",
        "render_targets"
      ],
      "state_space_snapshot": [
        "services_live",
        "service_interop_live",
        "env_live",
        "toolchain_live",
        "observability_live"
      ],
      "drift": [
        "missing_declared",
        "unexpected_live",
        "version_mismatch",
        "stale_state"
      ],
      "gate_inputs": [
        "evidence_refs",
        "freshness"
      ],
      "gate_decision": [
        "VALID",
        "INVALID",
        "UNKNOWN",
        "reasons[]",
        "blocking_edges[]"
      ]
    },
    "dependency_graph_model": {
      "nodes": [
        "service",
        "socket",
        "timer",
        "path",
        "slice",
        "command",
        "wrapper",
        "desktop_entry",
        "env_var",
        "config_file",
        "state_path",
        "credential",
        "package",
        "endpoint"
      ],
      "edges": [
        "depends_on",
        "wants",
        "activates",
        "launches",
        "wraps",
        "reads",
        "writes",
        "produces",
        "consumes",
        "listens_on",
        "exports_to",
        "conflicts_with"
      ],
      "provenance": [
        "origin",
        "evidence_ref",
        "confidence",
        "criticality"
      ]
    },
    "required_schemas": [
      "runtime.surface.schema.json",
      "userland.spec.schema.json",
      "state.space.schema.json",
      "dependency.graph.schema.json",
      "state.transition.request.schema.json",
      "state.transition.decision.schema.json"
    ]
  },
  "build_pipeline": {
    "passes": [
      "declared graph from runtime_surface[] + aggregate spec",
      "live graph from runtime (systemctl, sockets, process/runtime)",
      "journald overlay (live_event signals/edges)",
      "optional inference (confidence-scored, non-authoritative)"
    ],
    "outputs": [
      "state_space_snapshot",
      "dependency_graph",
      "drift"
    ]
  },
  "gate_semantics": {
    "decision_domain": [
      "VALID",
      "INVALID",
      "UNKNOWN"
    ],
    "rules": [
      "stale snapshot or missing hard-edge evidence => UNKNOWN",
      "inferred low-confidence edges must not hard-block by default",
      "gate output must carry rationale + evidence refs"
    ]
  },
  "lifecycle_and_triggers": {
    "meta_skill_surface": [
      "state-space build",
      "state-space update",
      "state-space gate",
      "state-space verify"
    ],
    "trigger_model": {
      "periodic": [
        "spawn-state-graph-refresh.timer",
        "Persistent=true",
        "RandomizedDelaySec"
      ],
      "post_change": [
        "deploy/config/unit/env/package/topology changes",
        "incremental first",
        "full rebuild fallback when confidence drops"
      ],
      "on_demand": [
        "spawnctl state refresh --full|--incremental"
      ]
    },
    "runtime_safety_controls": [
      "debounce change bursts",
      "single-writer lock for recomputes",
      "atomic publish (*.tmp -> rename)",
      "snapshot freshness required on transition checks"
    ]
  },
  "artifact_locations": {
    "current_graph": "$XDG_STATE_HOME/spawn/state/dependency.graph.current.json",
    "history": "$XDG_STATE_HOME/spawn/state/dependency.graph/<snapshot_id>.json",
    "required_fields": [
      "snapshot_id",
      "collected_at",
      "sources",
      "build_mode",
      "confidence_summary"
    ]
  },
  "scm_gate_integration": {
    "pre_commit": [
      "state-space verify",
      "state-space gate"
    ],
    "ci": [
      "state-space build in clean runner",
      "verify + gate as required checks"
    ],
    "post_commit": [
      "state-space update --scope incremental --reason post-commit",
      "queue full rebuild when confidence/coverage drops"
    ]
  },
  "adoption_sequence": [
    "build read-only state + graph artifacts",
    "run gates in shadow mode",
    "promote gates to enforcement once false positives are acceptable",
    "remove duplicated gate logic from PLAN/IMPLEMENT and keep LOOP orchestration-only"
  ]
}
